# contrib/systemd/backup-clean.service
[Unit]
Description=Periodic backup clean

[Service]
Type=oneshot
ExecStart=/home/marc/bin/backup-clean

[Install]
WantedBy=multi-user.target
